# python-authenticator.service
[Unit]
Description=Python Authenticator

[Service]
ExecStart=/usr/bin/env python /usr/local/bin/apiserver.py

[Install]
WantedBy=default.target
